# doc/script/lib/systemd/system/Triasinit.service
[Unit]
Description=tagent
After=network.target remote-fs.target nss-lookup.target

[Service]
Type=simple
ExecStart=/etc/init.d/Trias start
ExecStop=/etc/init.d/Trias stop
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target